# script/data/etc_init_datomic.conf
#! /bin/sh
### BEGIN INIT INFO
# Provides:           datomic
# Short-Description:  Start datomic transactor
### END INIT INFO
#
# Copyright (c) Vital Reactor LLC
#
# description: starts up datomic transactor in fixed location
env PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin:$PATH
env HOME=/opt/xn_apps
env DAEMON=/opt/datomic/current/bin/transactor
# If you want custom Xmx or Xms values, set them before the property file like so:
#env DAEMON_OPTS='-Xmx4g /etc/datomic/transactor.properties'
env DAEMON_OPTS='/etc/datomic/transactor.properties'
env DESC=DatomicTransactor
env APPDIR=/opt/xn_apps

env enabled=1

setuid root
setgid root

exec test -x $DAEMON || echo "$DESC Already running" && exit 1
description "$DESC"
 
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

pre-start script
        if [ ! -d "$APPDIR" ]; then
                mkdir -p $APPDIR
        fi
end script

script
        $DAEMON $DAEMON_OPTS >> /var/log/datomic.log 2>&1
end script
